# acs/lab5/slurm-node/supervisord.conf
[supervisord]
nodaemon=true

[program:munged]
command=/usr/sbin/munged -f --num-threads=1

[program:slurmctld]
command=/usr/sbin/slurmctld -D -vvv
